# package/DEBIAN/postinst
#!/usr/bin/env bash

set -e

systemctl daemon-reload && \
    systemctl start example-service && \
    systemctl enable example-service
