# package/DEBIAN/postinst
#!/usr/bin/env bash


set -ex

systemctl start example-service && \
    systemctl enable example-service
